# app/share/autostart/powerdwarf.desktop
[Desktop Entry]
Name=powerdwarf
Comment=Power Manager
Icon=battery
Exec=powerdwarf
Terminal=false
Type=Application
NotShowIn=GNOME;KDE;MATE;TDE;Unity;XFCE;Cinnamon;
StartupNotify=false
